# gpttg-bot.service
[Unit]
Description=GPTTG Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/GPTTG
ExecStartPre=/usr/bin/env mkdir -p /root/GPTTG/logs
# «апуск через Poetry, чтобы не зависеть от пути .venv
ExecStart=/usr/bin/env bash -c 'source /etc/profile && cd /root/GPTTG && poetry run python -m bot.main'
Restart=always
RestartSec=5
EnvironmentFile=/root/GPTTG/.env
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target