# gpttg-bot.service
[Unit]
Description=GPTTG Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/GPTTG
ExecStart=/root/GPTTG/.venv/bin/python3 -m bot.main
Restart=always
RestartSec=5
StandardOutput=append:/root/GPTTG/logs/bot.log
StandardError=append:/root/GPTTG/logs/bot.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
# Save this file as gpttg-bot.service in /etc/systemd/system/