# gpttg-bot.service
[Unit]
Description=GPTTG Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/GPTTG
ExecStartPre=/bin/mkdir -p /root/GPTTG/logs
ExecStart=/root/GPTTG/.venv/bin/python3 -m bot.main
Restart=always
RestartSec=5
StandardOutput=file:/root/GPTTG/logs/bot.log
StandardError=file:/root/GPTTG/logs/bot.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
# Save this file as gpttg-bot.service in /etc/systemd/system/