# bot/deploy/systemd/gpttg-update.timer
[Unit]
Description=Daily GPTTG update (10:00 UTC ≈ 13:00 MSK)

[Timer]
OnCalendar=*-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
