# bot/deploy/systemd/gpttg-update.timer
[Unit]
Description=Daily GPTTG update at 13:00 MSK (10:00 UTC)

[Timer]
OnCalendar=*-*-* 10:00:00      # 10:00 UTC ≈ 13:00 Moscow
Persistent=true

[Install]
WantedBy=timers.target
