# bot/deploy/systemd/gpttg-update.timer
[Unit]
Description=Daily GPTTG update at 13:00 MSK (10:00 UTC)

[Timer]
# 10:00 UTC ≈ 13:00 Europe/Moscow
OnCalendar=*-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target